# scripts/clawyard-queue-worker.conf
; ClawYard queue worker — Supervisor daemon
; 2026-05-28 v2: aumento de timeout 120→900s + tries 3→1 + stopwaitsecs 10→900.
;
; CAUSA: o user reportou "✗ Análise falhou após 95s — worker timeout" em
; tender 17662. Logs em prod mostravam MaxAttemptsExceededException
; recorrente. Root cause: este config tinha --timeout=120 que aborta jobs
; que excedem 2min, mas RunTenderAnalysisJob legítimo demora 60-180s (5-8
; agentes × 15-30s cada). Worker matava o job a meio.
;
; Adicionalmente --tries=3 sobrescrevia o `public int $tries = 1` do
; RunTenderAnalysisJob (intencional para evitar re-cobrar Anthropic
; em retries). Job era retentado 3× consumindo 3× Anthropic.
;
; stopwaitsecs=10 também era curto demais: durante deploys (Forge envia
; SIGTERM aos workers), 10s não chega para um panel multi-agente em curso.
; 900s = 15min, alinhado com o $timeout do job mais longo.
;
; APLICAR EM PROD (1× sudo, sem deploy):
;   scp scripts/clawyard-queue-worker.conf forge@clawyard.partyard.eu:/home/forge/
;   ssh forge@clawyard.partyard.eu
;   sudo cp /home/forge/clawyard-queue-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update clawyard-queue
;   sudo supervisorctl status clawyard-queue:*    # verifica RUNNING

[program:clawyard-queue]
directory=/home/forge/clawyard.partyard.eu/current/
command=php8.4 artisan queue:work --queue=high,default --tries=1 --timeout=900 --memory=256 --sleep=2 --max-time=3600
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
user=forge
startsecs=1
redirect_stderr=true
stdout_logfile=/home/forge/.forge/clawyard-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
; 900s = 15min — chega para o job mais longo ($timeout=600s) + buffer.
; Antes 10s matava panels multi-agente durante deploys.
stopwaitsecs=900
stopsignal=SIGTERM
stopasgroup=true
killasgroup=true
