# llm-proxy/supervisor/partyard-llm-proxy.conf
; PartYard LLM Proxy — supervisor config for Forge deployment.
;
; This file is NOT installed by deploy.sh (Forge's Daemons UI owns the
; canonical supervisor config). It's committed here so the config we ask
; Forge to create is auditable in git.
;
; To create the daemon in Forge:
;   Dashboard → Server → Daemons → New Daemon
;     Command:           /home/forge/llm-proxy/.venv/bin/uvicorn app:app --host 127.0.0.1 --port 8787 --workers 2 --timeout-keep-alive 75 --log-level warning --no-access-log --proxy-headers --forwarded-allow-ips=127.0.0.1
;     User:              forge
;     Directory:         /home/forge/llm-proxy
;     Processes:         1
;     Startsecs:         3
;     Stopwaitsecs:      30
;     Stopsignal:        TERM
;
; Forge writes the equivalent of the [program:xxx] block below to
; /etc/supervisor/conf.d/daemon-<ID>.conf and calls supervisorctl reread
; && supervisorctl update. The NOPASSWD sudo on /usr/bin/supervisorctl
; means we can restart without root password from deploy.sh.

[program:partyard-llm-proxy]
command=/home/forge/llm-proxy/.venv/bin/uvicorn app:app
    --host 127.0.0.1
    --port 8787
    --workers 2
    --timeout-keep-alive 75
    --log-level warning
    --no-access-log
    --proxy-headers
    --forwarded-allow-ips=127.0.0.1
user=forge
directory=/home/forge/llm-proxy
environment=PROXY_LOG_DIR="/home/forge/llm-proxy/logs",UPSTREAM_BASE_URL="https://api.anthropic.com"
autostart=true
autorestart=true
startsecs=3
startretries=10
stopwaitsecs=30
stopsignal=TERM
stdout_logfile=/home/forge/llm-proxy/logs/supervisor.stdout.log
stderr_logfile=/home/forge/llm-proxy/logs/supervisor.stderr.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
